# legacy/gluon-config-mode-ffgt/files/lib/gluon/setup-mode/notify-setup.cron
* * * * * /lib/gluon/setup-mode/notify-setup.sh
